# package/luci-app-sheepfold-family-internet-control/root/etc/init.d/sheepfold
#!/bin/sh /etc/rc.common

START=90
USE_PROCD=1

restore_usb_storage() {
        local enabled device role swap_file

        enabled="$(uci -q get sheepfold.usb.enabled 2>/dev/null || printf 0)"
        [ "$enabled" = 1 ] || return 0
        device="$(uci -q get sheepfold.usb.device 2>/dev/null || true)"
        [ -n "$device" ] || return 0

        if ! /usr/libexec/sheepfold/sheepfold-usb-storage mount "$device" >/dev/null 2>&1; then
                logger -t sheepfold "Configured USB device could not be mounted: $device"
                return 0
        fi

        role="$(uci -q get sheepfold.usb.role 2>/dev/null || printf logs_only)"
        [ "$role" = swap_logs ] || return 0
        swap_file='/mnt/sheepfold-usb/sheepfold/.swap'
        [ -f "$swap_file" ] || return 0
        grep -q "^$swap_file " /proc/swaps 2>/dev/null || \
                swapon "$swap_file" >/dev/null 2>&1 || \
                logger -t sheepfold "Configured Sheepfold swap file could not be enabled"
}

apply_led_repair() {
        [ -x /usr/libexec/sheepfold/sheepfold-router-control ] || return 0
        /usr/libexec/sheepfold/sheepfold-router-control led-repair-apply >/dev/null 2>&1 || \
                logger -t sheepfold "Configured router LED repair could not be applied"
}

apply_led_settings() {
        [ -x /usr/libexec/sheepfold/sheepfold-router-control ] || return 0
        /usr/libexec/sheepfold/sheepfold-router-control led-apply >/dev/null 2>&1 || \
                logger -t sheepfold "Configured router LED mode could not be applied"
}

start_service() {
        local port https_port

        restore_usb_storage
        apply_led_repair
        apply_led_settings

        [ -x /usr/libexec/sheepfold/sheepfold-site-lists ] &&
                /usr/libexec/sheepfold/sheepfold-site-lists cron-apply >/dev/null 2>&1 || true

        procd_open_instance
        procd_set_param command /usr/libexec/sheepfold/sheepfold-service
        procd_set_param respawn
        procd_close_instance

        if [ "$(uci -q get sheepfold.global.active_messenger 2>/dev/null || printf disabled)" = "telegram" ]; then
                procd_open_instance telegram
                procd_set_param command /usr/libexec/sheepfold/sheepfold-telegram-bot run
                procd_set_param respawn
                procd_close_instance
        fi

        port="$(uci -q get sheepfold.global.app_port 2>/dev/null || printf 5201)"
        https_port="$(uci -q get sheepfold.global.app_https_port 2>/dev/null || printf 5200)"
        case "$port" in ''|*[!0-9]*) port=5201 ;; esac
        case "$https_port" in ''|*[!0-9]*) https_port=5200 ;; esac

        procd_open_instance api
        if [ -s /etc/uhttpd.crt ] && [ -s /etc/uhttpd.key ] && \
           /usr/sbin/uhttpd -h 2>&1 | grep -q -- '-s'; then
                # HTTPS — основной транспорт, HTTP остаётся совместимым LAN fallback
                procd_set_param command /usr/sbin/uhttpd -f -h /www -x /cgi-bin \
                        -s "0.0.0.0:$https_port" -C /etc/uhttpd.crt -K /etc/uhttpd.key \
                        -p "0.0.0.0:$port" -T 30
        else
                logger -t sheepfold "TLS certificate or uhttpd TLS support is unavailable; using LAN HTTP fallback"
                procd_set_param command /usr/sbin/uhttpd -f -h /www -x /cgi-bin -p "0.0.0.0:$port" -T 30
        fi
        procd_set_param respawn
        procd_close_instance
}

reload_service() {
        /usr/libexec/sheepfold/sheepfold-service reload
        apply_led_repair
        apply_led_settings
}

stop_service() {
        [ -x /usr/libexec/sheepfold/sheepfold-router-control ] &&
                /usr/libexec/sheepfold/sheepfold-router-control blocked-page-stop >/dev/null 2>&1 || true
}
