# package/luci-app-sheepfold-family-internet-control/root/etc/cron.d/sheepfold-maintenance
# Тяжёлые операции в окне техобслуживания (02:00-05:00)

# Обновление пакета
30 2 * * * root /usr/libexec/sheepfold/sheepfold-maintenance-window && /usr/libexec/sheepfold/sheepfold-updater check >/dev/null 2>&1 || true

# Пересчёт ленты ARP/DHCP (дорогостоящая операция на слабых роутерах)
0 3 * * * root /usr/libexec/sheepfold/sheepfold-maintenance-window && /usr/libexec/sheepfold/sheepfold-device-detector scan full >/dev/null 2>&1 || true

# Глубокая синхронизация списков сайтов (site-lists)
0 4 * * * root /usr/libexec/sheepfold/sheepfold-maintenance-window && /usr/libexec/sheepfold/sheepfold-site-lists sync >/dev/null 2>&1 || true
